# ollama-cachyos/systemd/ollama.service
#============================================================================
# Title:            ollama.service
# Description:      Systemd service for Ollama Language Model Server with GPU optimizations
# Author:           Keverall
# Date:             2026-04-28
# Version:          1.0.0
# Documentation:    https://github.com/ollama/ollama
#============================================================================

[Unit]
Description=Ollama Language Model Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ollama
Group=ollama
# Set environment variables for optimization
Environment=OLLAMA_HOST=::
Environment=OLLAMA_PORT=11434
Environment=OLLAMA_NUM_PARALLEL=24
Environment=OLLAMA_MAX_LOADED_MODELS=2
# Optional: Set nice level for CPU priority (lower numbers = higher priority)
# Nice=10
# Optional: Set I/O scheduling class and priority
# IoniceClass=2
# IonicePriority=3
ExecStart=/usr/local/bin/ollama serve
ExecStop=/usr/bin/pkill -f "ollama serve"
Restart=always
RestartSec=5
# Limit memory usage (adjust based on your needs)
# MemoryLimit=20G
# Give the service access to GPU devices
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-uvm rw
DeviceAllow=/dev/nvidia-uvm-tools rw
# Standard output and error
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target