# ollama-devops/systemd/ollama.service
#============================================================================
# Title:            ollama.service
# Description:      Systemd service for Ollama Language Model Server
#                   Managed by sod.sh/eod.sh scripts with platform-aware config
# Author:           Keverall
# Date:             2026-04-30
# Version:          2.0.0
# Documentation:    https://github.com/ollama/ollama
#============================================================================

[Unit]
Description=Ollama Language Model Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ollama
Group=ollama

# Platform-aware configuration loaded via EnvironmentFile
# Default values provided; scripts override per-platform in /etc/default/ollama
EnvironmentFile=-/etc/default/ollama
EnvironmentFile=-/etc/sysconfig/ollama

# Core settings (fallbacks if env files missing)
Environment=OLLAMA_HOST=::
Environment=OLLAMA_PORT=11434
Environment=OLLAMA_NUM_PARALLEL=4
Environment=OLLAMA_MAX_LOADED_MODELS=1

# Conditional restart policy (not always-on; controlled by sod/eod scripts)
# Restart on failure OR when manually restarted via systemctl
Restart=ondemand
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Executable and arguments
ExecStart=/usr/local/bin/ollama serve
ExecStop=/bin/kill -TERM $MAINPID

# GPU device access - only active if OLLAMA_GPU_LAYERS is set > 0
# Platform-specific drop-in files can add DeviceAllow directives
# See: /etc/systemd/system/ollama.service.d/platform.conf

# Standard output and error (configurable)
StandardOutput=journal
StandardError=journal

# Systemd hardening defaults (adjust for GPU memory needs)
# MemoryLimit=20G
# NoNewPrivileges=no  # Required for NVIDIA driver access

[Install]
WantedBy=multi-user.target